# deploy/empmonitor-email.timer
[Unit]
Description=EMP Monitor — Nightly timesheet email at 20:00 SAST

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
